# dynv6.timer
[Unit]
Description=Runs dynv6 after boot

[Timer]
OnUnitActiveSec=10m
OnBootSec=1m
Unit=dynv6.service

[Install]
WantedBy=multi-user.target
