# dynv6-timer.service
[Unit]
Description=Dynv6 Timer

[Service]
ExecStart=/usr/local/bin/dynv6
